# meta-tanowrt/recipes-core/netifd/netifd/files/rootfs/etc/uci-defaults/15_migrate-network-config
#
# SPDX-License-Identifier: MIT
# 
# Copyright (C) 2021 Tano Systems LLC
# Anton Kikin <a.kikin@tano-systems.com>
#
/usr/libexec/netifd-utils/netifd-config-migrate.sh
