# meta-tanowrt/recipes-networking/vsftpd/vsftpd/files/vsftpd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=50
STOP=05

start() {
	/usr/sbin/vsftpd_prepare
	if [ "$?" = "0" ]; then
		service_start /usr/sbin/vsftpd_wrapper ipv4
		service_start /usr/sbin/vsftpd_wrapper ipv6
	fi
}

stop() {
	service_stop /usr/sbin/vsftpd
}
